# server/debian/efd-server.service
[Unit]
Description=efd-station SDR backend for Elad FDM-DUO
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/efd-server
Restart=always
RestartSec=2
Environment=RUST_LOG=info

# Run as dedicated user
User=efd
Group=efd

# USB and audio access
SupplementaryGroups=dialout audio plugdev

# Hardening
ProtectSystem=strict
ProtectHome=read-only
# Phase 2 persistence writes to $XDG_STATE_HOME (~/.local/state/efd-backend);
# config still lives under ~/.config/efd-backend. Both must be writable.
ReadWritePaths=/home/efd/.config/efd-backend /home/efd/.local/state/efd-backend
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
