# server/debian/efd-server.service
[Unit]
Description=efd-station SDR backend for Elad FDM-DUO
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/efd-server
Restart=always
# Short restart delay so client-initiated SelectDevice respawns feel
# snappy; process is designed to come up in <100 ms once systemd
# re-execs it.
RestartSec=500ms
# Bound the drain: `SHUTDOWN_TIMEOUT` in main.rs is 1 s, plus headroom
# for tokio + process exit. Don't sit 90 s (systemd default) on a
# stuck hardware task.
TimeoutStopSec=3
Environment=RUST_LOG=info

# Run as dedicated user
User=efd
Group=efd

# USB and audio access
SupplementaryGroups=dialout audio plugdev

# Hardening
ProtectSystem=strict
ProtectHome=read-only
# Phase 2 persistence writes to $XDG_STATE_HOME (~/.local/state/efd-backend);
# config still lives under ~/.config/efd-backend. Both must be writable.
# The `-` prefix on the state path tolerates it not existing yet — systemd
# fails to set up the mount namespace (status=226/NAMESPACE) otherwise,
# because ReadWritePaths bind-mounts the source directory. The path *is*
# created by postinst for fresh installs; the prefix guards upgrade paths
# where the old package version never had a chance to create it.
ReadWritePaths=/home/efd/.config/efd-backend -/home/efd/.local/state/efd-backend
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
