# server/debian/postinst
#!/bin/sh
set -e

# Create efd system user if it doesn't exist
if ! getent passwd efd > /dev/null 2>&1; then
    useradd --system --home-dir /home/efd --create-home \
        --groups dialout,audio,plugdev --shell /usr/sbin/nologin efd
fi

# Ensure config directory exists
mkdir -p /home/efd/.config/efd-backend
if [ ! -f /home/efd/.config/efd-backend/config.toml ]; then
    cp /etc/efd-backend/config.toml.example /home/efd/.config/efd-backend/config.toml
fi
chown -R efd:efd /home/efd/.config/efd-backend

# Ensure state directory exists (XDG_STATE_HOME — phase 2 persistence).
# Must exist at first service start so the hardened unit's ReadWritePaths
# can bind-mount it. We also keep a `-` prefix in the unit for upgrade
# safety, but creating it here makes first-boot behaviour clean.
mkdir -p /home/efd/.local/state/efd-backend
# Phase 4 REC: default recordings directory. Lives under state so
# persistence + recordings share one ReadWritePath. Config can point
# somewhere else; if so, that path needs its own writable slot.
mkdir -p /home/efd/.local/state/efd-backend/recordings
chown -R efd:efd /home/efd/.local

# Reload udev rules for FDM-DUO symlink
udevadm control --reload-rules 2>/dev/null || true
udevadm trigger 2>/dev/null || true

#DEBHELPER#

echo ""
echo "efd-server installed. To start:"
echo "  sudo systemctl enable --now efd-server"
echo ""
echo "Config: /home/efd/.config/efd-backend/config.toml"
echo ""
